# systemd/nanoclaw.service
[Unit]
Description=NanoClaw Personal Claude Assistant
After=network.target

[Service]
Type=simple
WorkingDirectory={{PROJECT_ROOT}}
ExecStart={{NODE_PATH}} {{PROJECT_ROOT}}/dist/index.js
Restart=always
RestartSec=5

Environment=PATH={{HOME}}/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME={{HOME}}
Environment=ASSISTANT_NAME=Andy

StandardOutput=append:{{PROJECT_ROOT}}/logs/nanoclaw.log
StandardError=append:{{PROJECT_ROOT}}/logs/nanoclaw.error.log

[Install]
WantedBy=default.target
